# config/openrc/fand-clientrpi.initd
#!/sbin/openrc-run
name="fand-clientrpi"
command="/usr/bin/env ${name}"
command_args="--logfile /var/log/${name}.log ${FAND_OPTS}"
start_stop_daemon_args="--respawn-delay 30 --respawn-max 120 --respawn-period 4000"
supervisor="supervise-daemon"
description="Simple daemon to control fan speed (${name})"

depend() {
	after net
}
